# sks-cron/etc/cron.d/cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root

0 0 * * * cd /var/sks && sudo -u sks sks cleandb
0 2 * * * sudo -u sks db_archive -d -h /var/sks/KDB && sudo -u sks db_archive -d -h /var/sks/PTree
