# sks-cron/etc/cron.d/cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root

0 0 * * * sks cleandb
0 2 * * * db_archive -d -h /var/lib/sks/KDB && db_archive -d -h /var/lib/sks/PTree
